# config/systemd/user/mopidy.service
[Unit]
Description=Mopidy

[Service]
ExecStart=/usr/bin/toolbox run mopidy
[Install]
WantedBy=default.target
